# deploy_tools/gunicorn-systemd.template.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/zw/sites/xfsu.top/gunicorn_blog/myblog.socket 

[Install]
WantedBy=sockets.target
